# packages/nspyre-drivers/nspyre-drivers-0.1.6.tar.gz/nspyre-drivers-0.1.6/src/nspyre_drivers/beaglebone/beaglebone_gpio/beaglebone_server/gpio_flask_server.service
[Unit]
Description=Remote control GPIO Flask server
After=network.target

[Service]
Type=idle
Restart=always
User=debian
ExecStart=/usr/bin/python3 /home/debian/beaglebone_server/gpio_flask_server.py

[Install]
WantedBy=multi-user.target
